# SkyeHands-main/AbovetheSkye-Platforms/skyeroutex-workforce-command-v0.4.0/deploy/skyeroutex-workforce-command.service
[Unit]
Description=SkyeRoutex Workforce Command
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skyeroutex
Group=skyeroutex
WorkingDirectory=/opt/skyeroutex-workforce-command
EnvironmentFile=/etc/skyeroutex-workforce-command/production.env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/skyeroutex-workforce-command

[Install]
WantedBy=multi-user.target
